# systemd/.config/systemd/user/backups.service
[Unit]
Description=Copy synced content to backup

[Service]
Type=simple
ExecStart=/run/media/evf/Storage/Backup/make_backups.sh
