# systemd/.config/systemd/user/paccleanup.service
[Unit]
Description=Pacman remove orphaned packages and clean cache

[Service]
Type=simple
ExecStart=/home/evf/scripts/paccleanup
